# deployments/podman/quadlet/hypersdk.container
# Podman Quadlet file for HyperSDK
# Place this file in /etc/containers/systemd/ or ~/.config/containers/systemd/
# Then run: systemctl --user daemon-reload && systemctl --user start hypersdk

[Unit]
Description=HyperSDK VM Export Daemon
Documentation=https://github.com/ssahani/hypersdk
After=network-online.target
Wants=network-online.target

[Container]
# Container image
Image=hypersdk/hypervisord:latest

# Container name
ContainerName=hypersdk-hypervisord

# Ports
PublishPort=8080:8080
PublishPort=8081:8081

# Volumes
Volume=hypersdk-data:/data:Z
Volume=hypersdk-exports:/exports:Z
Volume=/etc/hypersdk/config.yaml:/config/config.yaml:ro,Z

# Environment variables
Environment=DAEMON_ADDR=0.0.0.0:8080
Environment=LOG_LEVEL=info
Environment=DATABASE_PATH=/data/hypersdk.db

# vSphere credentials (can be overridden with EnvironmentFile)
# Environment=GOVC_URL=https://vcenter.example.com/sdk
# Environment=GOVC_USERNAME=administrator@vsphere.local
# Environment=GOVC_PASSWORD=changeme

# Use environment file for sensitive data
# EnvironmentFile=/etc/hypersdk/credentials.env

# Security
User=1000
Group=1000
SecurityLabelDisable=false

# Health check
HealthCmd=curl -f http://localhost:8080/health || exit 1
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3

# Networking
Network=hypersdk-network

# Restart policy
Restart=always

# Logging
LogDriver=journald

# Resource limits
MemoryLimit=2G
CPUQuota=100%

[Service]
# Service configuration
Restart=always
RestartSec=10s
TimeoutStartSec=900

# Ensure volumes are created
ExecStartPre=/usr/bin/podman volume create hypersdk-data
ExecStartPre=/usr/bin/podman volume create hypersdk-exports

# Ensure network exists
ExecStartPre=/usr/bin/podman network exists hypersdk-network || /usr/bin/podman network create hypersdk-network

[Install]
WantedBy=multi-user.target default.target
