# dist/hypersdk-0.2.0-linux-amd64/docs/hypervisord.service
[Unit]
Description=HyperSDK VM Export Daemon
After=network.target

[Service]
Type=simple
User=hypersdk
Group=hypersdk
WorkingDirectory=/opt/hypersdk
ExecStart=/opt/hypersdk/bin/hypervisord --config /etc/hypersdk/config.yaml
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hypervisord

[Install]
WantedBy=multi-user.target
